# scripts/sv/hostapd/run
#!/usr/bin/bash

exec chpst -u root /bin/hostapd

