# scripts/station/sv/wpa/run
#!/bin/sh -e

while [ ! -r "/sys/class/net/wlan0/address" ]; do
    sleep 2
done

exec 2>&1
exec chpst -u root /usr/sbin/wpa_supplicant -i wlan0 -c /etc/wpa_supplicant/wpa_supplicant.conf

